# s6/services.d/reproxy/finish
#!/usr/bin/execlineb -S0

# if reproxy fails - container fails too
s6-svscanctl -t /var/run/s6/services
