# s6/services.d/backend/finish
#!/usr/bin/execlineb -S0

# if backend fails - container fails too
s6-svscanctl -t /var/run/s6/services
