# s6/services.d/backend/run
#!/usr/bin/with-contenv ash

exec /usr/bin/backend
